# rserver.service.template
[Unit]
Description=rserver

[Service]
WorkingDirectory=PWD
ExecStart=Rscript server/runner.R
